# scripts/bsd-controller.service
[Unit]
Description=BSD Unitree Controller (FastAPI + ROS 2 节点)
After=network.target

[Service]
# 用 unitree 用户运行，不用 root
User=unitree
Group=unitree

# 工作目录
WorkingDirectory=/home/unitree/bsd-unitree-controller

# 启动命令：调用 start.sh
ExecStart=/home/unitree/bsd-unitree-controller/scripts/start.sh

# 自动重启：进程异常退出后 5 秒重启
Restart=always
RestartSec=5

# 日志输出到 systemd journal，用 journalctl 查看
StandardOutput=journal
StandardError=journal

# 环境变量（HOME 让 pip --user 装的包能被找到）
Environment=HOME=/home/unitree
Environment=PYTHONUNBUFFERED=1

[Install]
# 开机自启
WantedBy=multi-user.target
